# .config/systemd/user/anybar.timer
[Unit]
Description=Anybar

[Timer]
OnUnitInactiveSec=30

[Install]
WantedBy=timers.target
