# .config/systemd/user/r0covid.timer
[Unit]
Description=R0covid

[Timer]
OnCalendar=*-*-* 20:01:00
Persistent=true

[Install]
WantedBy=timers.target
