# .config/systemd/user/smscheck.timer
[Unit]
Description=Check new SMS

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
