# .config/systemd/user/kunst.service.d/override.conf
[Service]
Environment="KUNST_MODE=notify"
